# dot_config/systemd/user/blueman-applet.service
[Unit]
Description=Bluetooth management applet
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/blueman-applet

[Install]
WantedBy=graphical-session.target
