# dot_config/systemd/user/xdg-desktop-portal-hyprland.service.d/override.conf
[Unit]
PartOf=hyprland.target
After=hyprland.target
