# dot_config/systemd/user/passh-agent.service.d/override.conf
[Unit]
After=gpg-agent.service

[Service]
ExecStartPre=/bin/sleep 6
